# systemd/gopro-sync.service
[Unit]
Description=GoPro Gardening hourly sync
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=__ROOT__
ExecStart=/usr/bin/env bash __ROOT__/scripts/run_sync.sh
